# illumos/pivy-fs.xml.in
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type='manifest' name='SUNWcsr:filesystem-pivy'>

<service
        name='system/filesystem/pivy'
        type='service'
        version='1'>

        <create_default_instance enabled='false' />

        <single_instance/>

        <dependency
                name='single-user'
                grouping='require_all'
                restart_on='none'
                type='service'>
                <service_fmri value='svc:/milestone/single-user' />
        </dependency>

        <dependency
                name='pivy-agent'
                grouping='require_all'
                restart_on='none'
                type='service'>
                <service_fmri value='svc:/system/pivy-agent' />
        </dependency>

        <dependent
                name='fs-local'
                grouping='require_all'
                restart_on='none'>
                <service_fmri value='svc:/system/filesystem/local' />
        </dependent>
        <exec_method
                type='method'
                name='start'
                exec='@@METHODPATH@@/fs-pivy'
                timeout_seconds='0' />

        <exec_method
                type='method'
                name='stop'
                exec=':true'
                timeout_seconds='0' />

        <property_group name='startd' type='framework'>
                <propval name='duration' type='astring' value='transient' />
        </property_group>

        <property_group name='config' type='application'>
        </property_group>

        <stability value='Unstable' />

        <template>
                <common_name>
                        <loctext xml:lang='C'>
                                Unlock encrypted ZFS datasets using pivy-zfs
                        </loctext>
                </common_name>
        </template>
</service>

</service_bundle>
